# systemd/fauxnet.target
[Unit]
Description=Fauxnet - Web Frontend for Greybox/Core
Requires=fauxnet@webui.service fauxnet@api.service

[Install]
WantedBy=multi-user.target